# infra/freebsd/rc.d/grappa
#!/bin/sh
#
# PROVIDE: grappa
# REQUIRE: NETWORKING SERVERS
# KEYWORD: shutdown
#
# Grappa IRC bouncer — FreeBSD rc.d wrapper for the `mix release` build.
#
# Adds to /etc/rc.conf (or /etc/rc.conf.d/grappa):
#   grappa_enable="YES"
#   grappa_user="grappa"                                   # default
#   grappa_release_root="/home/grappa/grappa/_build/prod/rel/grappa"
#   grappa_env_file="/usr/local/etc/grappa/grappa.env"     # optional
#
# The env file (chmod 600, owned by grappa_user) sources secrets:
#   SECRET_KEY_BASE=...
#   SECRET_SIGNING_SALT=...
#   RELEASE_COOKIE=...
#   GRAPPA_ENCRYPTION_KEY=...
#   VAPID_PUBLIC_KEY=...
#   VAPID_PRIVATE_KEY=...
#   GRAPPA_OUTBOUND_V6_POOL=2a03:4000:2:33c::9000,2a03:4000:2:33c::442,...
#   DATABASE_PATH=/home/grappa/grappa/runtime/grappa_prod.db
#   PHX_HOST=grappa.bad.ass
#   PORT=4000
#
# Why a thin wrapper around `bin/grappa daemon` instead of `daemon(8)`:
# the release script already daemonizes + writes its own pid file, and
# routes BEAM signals correctly. Reinventing that with the host
# daemon(8) loses Erlang's signal handling (SIGTERM → graceful
# Application.stop callback chain → QUIT to upstream peers).

. /etc/rc.subr

name="grappa"
rcvar="grappa_enable"

load_rc_config "$name"

: ${grappa_enable:="NO"}
: ${grappa_user:="grappa"}
: ${grappa_release_root:="/home/grappa/grappa/_build/prod/rel/grappa"}
: ${grappa_env_file:=""}
# Where run_erl (mix release daemon) writes its `log/erlang.log.*`,
# `pipe/`, and `pid` files. Defaults to `<repo>/runtime` so the
# erlang.log directory lands at the canonical `runtime/log/` path
# (sibling of `runtime/grappa_prod.db`, `runtime/uploads/`). Setting
# this to `runtime/log` would produce `runtime/log/log/erlang.log.*`
# because run_erl always creates its own `log/` subdir under
# RELEASE_TMP. Operator override via /etc/rc.conf.d/grappa.
: ${grappa_runtime_tmp:="/home/grappa/grappa/runtime"}
# Stop/start synchronization (defect #9, 2026-06-11 outage: `service
# grappa restart` started the new BEAM while the old node was still
# draining WS connections → "name grappa@grappa seems to be in use" →
# new BEAM died at boot and rc.d printed "Starting grappa." and walked
# away). grappa_stop blocks until BEAM exit + epmd name release;
# grappa_start refuses a doomed boot and verifies the node comes up.
# The wait logic is shared with deploy.sh's cold path — one
# implementation, jail_beam_wait.sh.
: ${grappa_beam_wait:="/home/grappa/grappa/infra/freebsd/jail_beam_wait.sh"}
: ${grappa_node:="grappa"}
: ${grappa_stop_timeout:="30"}
: ${grappa_start_timeout:="30"}
: ${grappa_name_wait_timeout:="15"}

command="${grappa_release_root}/bin/grappa"
pidfile="${grappa_runtime_tmp}/pid"
extra_commands="reload remote_console release_info"

start_cmd="grappa_start"
stop_cmd="grappa_stop"
status_cmd="grappa_status"
remote_console_cmd="grappa_remote_console"
release_info_cmd="grappa_release_info"

grappa_env_export() {
	# Source the env file as the grappa user — keeps secrets readable
	# only by grappa_user (chmod 600). `set -a` exports every var the
	# file defines so the release inherits them.
	if [ -n "${grappa_env_file}" ] && [ -r "${grappa_env_file}" ]; then
		echo "[grappa] sourcing env from ${grappa_env_file}"
	elif [ -n "${grappa_env_file}" ]; then
		echo "[grappa] WARNING: env file ${grappa_env_file} not readable by current user — release will boot with empty env (likely to crash on missing SECRET_KEY_BASE)" >&2
	fi
}

grappa_runas() {
	# Run a release subcommand as grappa_user with the env file sourced.
	# The env file (if set) is read inside the subshell so secrets are
	# only in memory for the lifetime of the child process.
	# RELEASE_TMP redirects run_erl's erlang.log.* + pid file from the
	# release's `_build/.../tmp/` (rebuilt on every `mix release`) to
	# the persistent runtime dir under `<repo>/runtime/log`. Must be
	# `export`-ed (not `VAR=val cmd` prefix) so it survives past the
	# `. ${grappa_env_file}` source and lands on the final ${command}.
	local subcmd="$1"
	local envsrc=""
	if [ -n "${grappa_env_file}" ] && [ -r "${grappa_env_file}" ]; then
		envsrc=". ${grappa_env_file} && export \$(grep -Eo '^[A-Z_][A-Z0-9_]*' ${grappa_env_file}) && "
	fi
	# rc(8) services inherit rc.subr's stock PATH, which has NO
	# /usr/local — pkg-installed binaries (exiftool, ffmpeg for
	# Grappa.Uploads.MetadataStrip #39) are invisible to the BEAM's
	# System.find_executable/1 without this prepend. Found live
	# 2026-06-10: pkgs installed, yet every media upload 422'd
	# fail-closed with "exiftool not found on PATH".
	su -m "${grappa_user}" -c "export RELEASE_TMP='${grappa_runtime_tmp}'; export PATH=\"/usr/local/bin:/usr/local/sbin:\${PATH}\"; ${envsrc}${command} ${subcmd}"
}

grappa_start() {
	grappa_env_export
	mkdir -p "${grappa_runtime_tmp}"
	chown "${grappa_user}":"${grappa_user}" "${grappa_runtime_tmp}"

	# Refuse to boot into a name collision. A registered name here
	# means an old node is still draining (e.g. a restart racing an
	# external stop) — booting now is guaranteed death at distribution
	# start. NO escalation on this side: the name's owner may be a
	# live BEAM that must not be shot (see jail_beam_wait.sh header).
	if [ -x "${grappa_beam_wait}" ]; then
		if ! "${grappa_beam_wait}" wait-name-free "${grappa_node}" "${grappa_name_wait_timeout}"; then
			echo "ERROR: refusing to start — epmd name '${grappa_node}' is still registered (old node draining?). Retry once \`epmd -names\` is clean." >&2
			return 1
		fi
	else
		echo "WARNING: ${grappa_beam_wait} missing — starting without the epmd name-collision guard" >&2
	fi

	echo "Starting grappa."
	grappa_runas "daemon"

	# `bin/grappa daemon` returns right after spawning run_erl; the
	# 2026-06-11 outage signature was "Starting grappa." as the last
	# line while the BEAM had already died on name-in-use. Verify the
	# node actually comes up: poll the release's `pid` RPC (the same
	# probe grappa_status uses), and treat a vanished beam.smp as an
	# immediate loud failure instead of waiting out the timeout.
	_seen_beam=0
	_i=0
	while :; do
		if grappa_runas "pid" >/dev/null 2>&1; then
			echo "grappa is up (node responsive)."
			return 0
		fi
		if pgrep -q beam.smp 2>/dev/null; then
			_seen_beam=1
		elif [ "${_seen_beam}" -eq 1 ]; then
			echo "ERROR: BEAM exited during boot — see ${grappa_runtime_tmp}/log/erlang.log.*" >&2
			return 1
		fi
		if [ "${_i}" -ge "${grappa_start_timeout}" ]; then
			echo "ERROR: node not responsive after ${grappa_start_timeout}s — see ${grappa_runtime_tmp}/log/erlang.log.*" >&2
			return 1
		fi
		_i=$((_i + 1))
		sleep 1
	done
}

grappa_stop() {
	echo "Stopping grappa."
	grappa_runas "stop"
	# `bin/grappa stop` is ASYNCHRONOUS — it signals the node and
	# returns while the BEAM is still draining WS connections and
	# running terminate callbacks. A `service grappa restart` that
	# trusts an async stop boots the new node into the old node's epmd
	# registration (defect #9 — THE 2026-06-11 outage). Block until
	# stop means STOPPED: BEAM exited, name released.
	if [ -x "${grappa_beam_wait}" ]; then
		"${grappa_beam_wait}" wait-stopped "${grappa_node}" "${grappa_stop_timeout}"
	else
		echo "WARNING: ${grappa_beam_wait} missing — stop returns while the BEAM may still be draining (a restart may race the epmd name)" >&2
	fi
}

grappa_status() {
	grappa_runas "pid" 2>/dev/null && echo "grappa is running." || echo "grappa is stopped."
}

grappa_remote_console() {
	grappa_runas "remote"
}

grappa_release_info() {
	grappa_runas "version"
}

run_rc_command "$1"
